# .config/autostart/brave.desktop
[Desktop Entry]
Type=Application
Name=Brave Autostart
Exec=flatpak run com.brave.Browser
Comment=Start Brave on login
